# tools/pi_agent/pi_agent.service
[Unit]
Description=Pi test agent for router passthrough testing
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/pi_agent/pi_agent.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
